# com.ym.fcb.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent 模板：feishu-claude-bridge
  放到 ~/Library/LaunchAgents/com.ym.fcb.plist 并修改为你自己的路径/凭证。

  关键配置:
    - RunAtLoad: 登录自启
    - KeepAlive.SuccessfulExit=false: 非零退出码 = 拉起；正常退出 = 不拉起
    - KeepAlive.NetworkState=true: 网络可用时才启动（合盖睡醒后断网不会疯转）
    - StandardOut/Err: bridge 的 log 落盘到 state dir

  ⚠️ 安全:
    plist 里存 ANTHROPIC_API_KEY 属于 plaintext-at-rest，同机任何能读 ~/Library/
    的 process 都能看到。如果担心:
      - 从 plist 里移除 ANTHROPIC_API_KEY
      - 用 `launchctl setenv ANTHROPIC_API_KEY sk-xxx` 手动注入 launchd 环境
      - 或者用 macOS Keychain + 启动脚本读出来再 exec bridge

  加载:
    launchctl load ~/Library/LaunchAgents/com.ym.fcb.plist
    launchctl start com.ym.fcb
    tail -f ~/.local/state/fcb/bridge.log

  卸载:
    launchctl unload ~/Library/LaunchAgents/com.ym.fcb.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ym.fcb</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/bun</string>
    <string>run</string>
    <string>/Users/ym/Projects/feishu-claude-bridge/bridge.ts</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/ym/Projects/feishu-claude-bridge</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>NetworkState</key>
    <true/>
  </dict>

  <!-- Throttle: 崩溃后至少等 10 秒才重启（防止 tight crash loop） -->
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/Users/ym/.local/state/fcb/bridge.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/ym/.local/state/fcb/bridge.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- 改成你自己的 open_id（拿自己的 id 见 README） -->
    <key>FCB_OWNER_OPEN_IDS</key>
    <string>ou_REPLACE_WITH_YOUR_OPEN_ID</string>

    <!-- Anthropic API key，也可改用 launchctl setenv（见顶部注释） -->
    <key>ANTHROPIC_API_KEY</key>
    <string>sk-ant-REPLACE_ME</string>

    <!-- 必须显式给 PATH，launchd 默认 PATH 很短 -->
    <key>PATH</key>
    <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/bin:/bin:/usr/sbin:/sbin</string>

    <!-- 可选：bridge 状态目录（不写默认 ~/.local/state/fcb） -->
    <!-- <key>FCB_STATE_DIR</key><string>/Users/ym/.local/state/fcb</string> -->
  </dict>
</dict>
</plist>
